# files/etc/crontabs/root
# Compact zram every 2 hours
0 */2 * * *	/bin/zram compact
